# deploy/supervisor/pathfy-backend.conf
; 复制到服务器: /etc/supervisor/conf.d/pathfy-backend.conf
; 然后: supervisorctl reread && supervisorctl update && supervisorctl start pathfy-backend
;
; 宝塔也可在「Supervisor 管理器」里新建进程，把下方 command / directory 填进去即可。

[program:pathfy-backend]
command=/www/wwwroot/pathfy-uni/backend/.venv/bin/gunicorn -w 2 -b 127.0.0.1:5000 "app:create_app()" --timeout 120 --access-logfile /www/wwwlogs/pathfy-backend-access.log --error-logfile /www/wwwlogs/pathfy-backend-error.log
directory=/www/wwwroot/pathfy-uni/backend
user=www
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=30
stopasgroup=true
killasgroup=true
priority=999
environment=PATH="/www/wwwroot/pathfy-uni/backend/.venv/bin"
